# resources/etc/init/ipchangecheck.conf
description "IP address change check"
author      "Robert Auer <robert.auer@cloudogu.com"

start on started etcd
stop on runlevel [!2345]

exec /bin/bash -c '/usr/local/bin/ipchange.sh'
