# resources/etc/init/ipchangecheck.conf
description "IP address change check"

start on (started etcd)
stop on runlevel [!2345]

exec /etc/ipchange.sh
