# resources/etc/systemd/system/docker.service.d/dockeroptions.conf
[Service]
Type=simple
ExecStart=
ExecStart=/bin/bash -c '/etc/systemd/system/docker.service.d/docker-start-script.sh'
Restart=always
