# resources/etc/systemd/system/docker.service.d/dockeroptions.conf
[Service]
ExecStart=/etc/systemd/system/docker.service.d/docker-start-script.sh
Restart=always
