# resources/etc/systemd/system/docker-metadata.service
[Service]
ExecStart=
ExecStart=/etc/systemd/system/docker-metadata.service.d/docker-start-script.sh
Restart=no
