# containers/registrator/resources/etc/init/ces-registrator.conf
description "CES registrator container"
author "Sebastian Sdorra <sebastian.sdorra@triology.de>"
start on started docker
stop on stopped docker
respawn
script
	/usr/bin/docker start -a registrator
end script
